# daml/infra/macos/scripts/com.digitalasset.macos-builder-crontab.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.digitalasset.macos-builder-crontab</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/builder/rebuild-crontask.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>20</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>AbandonProcessGroup</key>
  <true/>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardErrorPath</key>
  <string>/Users/builder/rebuild-crontab.err</string>

  <key>StandardOutPath</key>
  <string>/Users/builder/rebuild-crontab.out</string>
</dict>
</plist>
